# bin/blackbox-backup.timer
[Unit]
Description=Run BlackBox backup daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
